# crontab
*/1 * * * * python manage.py tradeapp.simulate.simulate_profit_loss
